# monitoring/co2/local/cron/cron.hourly/rrdbackup
#!/bin/sh

RRDDATA=/dev/shm/rrd.db
RRDBACKUP=/var/lib/rrd/co2

rsync -avx --delete $RRDDATA/ $RRDBACKUP/